# mpcontribs-api/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:worker]
command=wait-for-it.sh %(ENV_JUPYTER_GATEWAY_HOST)s -s -t 50 -- wait-for-it.sh %(ENV_MPCONTRIBS_API_HOST)s -s -t 15 -- flask rq worker
numprocs=1
directory=/app
autostart=false
autorestart=false
stopsignal=TERM
exitcodes=124
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:scheduler]
command=wait-for-it.sh %(ENV_JUPYTER_GATEWAY_HOST)s -s -t 50 -- wait-for-it.sh %(ENV_MPCONTRIBS_API_HOST)s -s -t 15 -- flask rq scheduler
numprocs=1
directory=/app
autostart=false
autorestart=false
stopsignal=TERM
exitcodes=124
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:api]
command=wait-for-it.sh %(ENV_JUPYTER_GATEWAY_HOST)s -t 50 -- ./start.sh
directory=/app
autostart=false
autorestart=false
stopsignal=TERM
exitcodes=124
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:main]
command=./main.sh
directory=/app
environment=METADATA_URI=%(ENV_ECS_CONTAINER_METADATA_URI_V4)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
